# firmware/X8SE/CM1C_FC-Relay/_app1553831646872CM1C-D-V010SP17RTM190327R17-CN-def_A-Release-9F4B6DDF6261B426340F8F5E2B7DF4A1.bin1553831646912.extracted/squashfs-root/etc/hotplug.d/iface/60-dnsmasq
#!/bin/sh

. /lib/functions.sh
. /lib/functions/network.sh

#restart dhcp server when domain changed and this domain should be delegated to LAN side
option_cb(){
	[ "$1" != "update_domain" ] && return
	[ "$2" != "$INTERFACE" ] && return

	old_domain=$(grep domain= /var/etc/dnsmasq.conf | cut -d "=" -f 2)
	network_get_dnssearch new_domain "$INTERFACE"
	[ "$old_domain" != "$new_domain" ] && /etc/init.d/dnsmasq restart
}

/etc/init.d/dnsmasq enabled && [ "$ACTION" = "ifup" ] && {
	config_load dhcp
}
