# firmware/X8SE/CM1C_FC-Relay/_app1553831646872CM1C-D-V010SP17RTM190327R17-CN-def_A-Release-9F4B6DDF6261B426340F8F5E2B7DF4A1.bin1553831646912.extracted/squashfs-root/etc/init.d/fmlink
#!/bin/sh /etc/rc.common

START=88

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1

start() {
	/usr/bin/fmlink_drone &
	sleep 5
	insmod /lib/modules/3.3.8/ath_udc.ko
	insmod /lib/modules/3.3.8/g_ether.ko
}

stop() {
	killall -9 fmlink_drone
}

