# firmware/X8SE/CM1C_FC-Relay/_app1553831646872CM1C-D-V010SP17RTM190327R17-CN-def_A-Release-9F4B6DDF6261B426340F8F5E2B7DF4A1.bin1553831646912.extracted/squashfs-root/etc/uci-defaults/sound
#!/bin/sh
#
# Copyright (c) 2013 The Linux Foundation. All rights reserved.
# Copyright (C) 2011 OpenWrt.org
#

touch /etc/config/sound

. /lib/functions/uci-defaults.sh
. /lib/ar71xx.sh

board=$(ar71xx_board_name)

case "$board" in
cus227*)
	ucidef_set_snd_ctrl "card0" "Right Mixer" "unmute"
	ucidef_set_snd_ctrl "card0" "Left Mixer" "unmute"
	;;
esac

uci commit sound

exit 0
