# firmware/X8SE/CM1C_FC-Relay/_app1553831646872CM1C-D-V010SP17RTM190327R17-CN-def_A-Release-9F4B6DDF6261B426340F8F5E2B7DF4A1.bin1553831646912.extracted/squashfs-root/usr/lib/opkg/info/base-files.list
/etc/openwrt_version
/sbin/sound
/lib/upgrade/openmesh.sh
/lib/ar71xx.sh
/etc/resolv.conf
/lib/upgrade/dir825.sh
/etc/init.d/done
/var
/lib/firstboot/99_10_no_fo_cleanup
/etc/hotplug2-init.rules
/etc/preinit
/lib/firstboot/50_pivot
/sbin/mount_root
/lib/firstboot/30_reset_copy_rom
/sbin/wget2nand
/lib/upgrade/platform.sh
/etc/hotplug2-common.rules
/etc/shells
/lib/preinit/05_set_preinit_iface_ar71xx
/lib/preinit/70_initramfs_test
/lib/preinit/50_choose_console
/lib/firstboot/10_no_fo_clear_overlay
/etc/uci-defaults/wrt160nl
/etc/protocols
/lib/preinit/03_preinit_do_ar71xx.sh
/lib/preinit/20_check_jffs2_ready
/lib/preinit/10_indicate_preinit
/lib/firstboot/99_10_with_fo_cleanup
/etc/hotplug.d/button/00-button
/etc/hotplug.d/net/10-ar922x-led-fix
/etc/config/system
/lib/preinit/40_mount_ubifs
/etc/fstab
/lib/preinit/30_device_fs_daemons
/lib/preinit/90_mount_no_jffs2
/etc/sysctl.conf
/sbin/sysupgrade
/lib/preinit/70_pivot_jffs2_root
/lib/preinit/99_10_run_init
/lib/firstboot/20_has_mini_fo
/lib/firstboot/10_reset_has_mini_fo
/lib/upgrade/ap135.sh
/lib/firstboot/10_determine_parts
/sbin/led.sh
/lib/preinit/50_indicate_regular_preinit
/lib/firstboot/40_no_fo_copy_ramoverlay
/etc/banner
/lib/preinit/10_essential_fs
/etc/shadow
/etc/rc.common
/lib/preinit/40_mount_ext4
/lib/preinit/10_indicate_failsafe
/lib/preinit/40_init_shm
/init
/lib/debug/system
/etc/init.d/watchdog
/lib/preinit/40_mount_jffs2
/lib/preinit/20_device_fs_mount
/lib/preinit/80_mount_root
/etc/rc.local
/lib/functions/uci-defaults.sh
/etc/init.d/defconfig
/etc/init.d/led
/lib/functions/network.sh
/etc/init.d/sysctl
/lib/firstboot/30_is_rootfs_mounted
/sbin/hotplug-call
/etc/uci-defaults/sound
/etc/init.d/usb
/lib/preinit/05_set_iface_mac_ar71xx
/lib/firstboot/20_format_mount_ext4
/lib/preinit/40_mount_devpts
/lib/preinit/90_init_console
/etc/services
/lib/preinit/99_10_mount_no_mtd
/lib/preinit/10_check_for_mtd
/lib/upgrade/common.sh
/etc/sysupgrade.conf
/lib/firstboot/05_firstboot_skip
/lib/upgrade/keep.d/base-files
/etc/uci-defaults/vlan-migration
/lib/preinit/40_run_failsafe_hook
/etc/group
/etc/uci-defaults/migrate-sysctl
/lib/preinit/30_failsafe_wait
/lib/preinit/05_enable_reset_button_ar71xx
/etc/functions.sh
/sbin/firstboot
/lib/upgrade/allnet.sh
/bin/login.sh
/etc/hosts
/etc/uci-defaults/migrate-shadow
/etc/defconfig/wndr3700/network
/etc/hotplug.d/usb/10-usb
/lib/functions.sh
/lib/preinit/60_init_hotplug
/lib/preinit/70_pivot_ext4_root
/sbin/sysdebug
/etc/init.d/rcS
/etc/uci-defaults/leds
/rom/note
/etc/inittab
/etc/openwrt_release
/lib/preinit/05_mount_skip
/sbin/wifi
/etc/profile
/lib/preinit/70_pivot_ubifs_root
/lib/preinit/81_load_wifi_board_bin
/lib/preinit/99_10_failsafe_login
/etc/TZ
/lib/firstboot/20_no_fo_mount_jffs
/etc/uci-defaults/inittab-console-fixup
/lib/functions/service.sh
/etc/init.d/umount
/lib/upgrade/cus227.sh
/lib/preinit/90_restore_config
/lib/upgrade/keep.d/base-files-essential
/lib/firstboot/30_no_fo_pivot
/etc/passwd
/etc/hotplug.d/ieee1394/10-ieee1394
/lib/preinit/02_default_set_state
/etc/mtab
/lib/firstboot/40_copy_ramoverlay
/etc/init.d/boot
/lib/functions/gpio-button-hotplug.sh
/lib/functions/boot.sh
/bin/ipcalc.sh
/etc/diag.sh
/lib/firstboot/20_reset_clear_jffs
/lib/preinit/41_merge_overlay_hooks
/etc/uci-defaults/network
