# firmware/X8SE/CM1C_FC-Relay/_app1553831646872CM1C-D-V010SP17RTM190327R17-CN-def_A-Release-9F4B6DDF6261B426340F8F5E2B7DF4A1.bin1553831646912.extracted/squashfs-root/usr/lib/opkg/info/netifd.list
/lib/netifd/proto/dhcp.sh
/etc/init.d/network
/lib/netifd/netifd-proto.sh
/etc/hotplug.d/iface/00-netstate
/sbin/ifdown
/lib/netifd/dhcp.script
/usr/share/udhcpc/default.script
/lib/network/config.sh
/etc/hotplug.d/iface/10-sysctl
/sbin/netifd
/sbin/devstatus
/sbin/ifstatus
/sbin/ifup
