# firmware/x6_devel/_X6E-RC-V050SP28RTM170621R104-CN-AA-Release-F9DD9F1FE50958DEBEF1055DA4110CF5.bin15135656614721513585575123.extracted/_258040.extracted/cpio-root/etc/inittab
::sysinit:/etc/rc.d/rcS
::respawn:/sbin/getty ttyS0 230400
::shutdown:/bin/umount -a

